# systemd/bitunix-recenter.timer
[Unit]
Description=Daily 08:00 BitUnix Grid Re-center
[Timer]
OnCalendar=*-*-* 08:00:00
Persistent=true
[Install]
WantedBy=timers.target
